# ueventmonitor-1.0/data/ueventmonitor.service.in
[Unit]
Description=uevent monitor

[Service]
Type=dbus
BusName=org.freedesktop.UeventMonitor
ExecStart=@ueventmonitorprivdir@/ueventmonitord
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target